# playbooks/roles/common-deployer/templates/helm-serve.service.j2
[Unit]
Description=Helm Server
After=network.target

[Service]
User={{ ansible_user }}
Restart=always
ExecStart=/usr/bin/helm serve

[Install]
WantedBy=multi-user.target
